# systemd/voms@.service.d/after-mariadb.conf
## Ensures that VOMS starts after MySQL db is loaded
[Unit]
After=mariadb.service
